# deploy/linux/phd-atlas.service
[Unit]
Description=PhD Atlas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phd-atlas
Group=phd-atlas
WorkingDirectory=/opt/phd-atlas
EnvironmentFile=/etc/phd-atlas/phd-atlas.env
ExecStart=/usr/bin/node tools/start-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=35
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/phd-atlas/storage

[Install]
WantedBy=multi-user.target
