# deploy/linux/phd-atlas.service
[Unit]
Description=PhD Atlas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phd-atlas
Group=phd-atlas
WorkingDirectory=/opt/phd-atlas
EnvironmentFile=/etc/phd-atlas/phd-atlas.env
ExecStart=/usr/bin/node tools/start-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=35
KillSignal=SIGTERM
# The short-lived detached update helper must survive the main process exit.
# start-server.mjs waits on its durable lock before the supervisor relaunches.
KillMode=process
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Admin update packages atomically replace the managed runtime and node_modules;
# application data remains under storage/.
ReadWritePaths=/opt/phd-atlas

[Install]
WantedBy=multi-user.target
